# root/etc/init/proxy.conf.erb
start on started nginx
respawn

expect daemon

<% vars.each do |name, value| %>
env <%= name %>='<%= value %>'
<% end %>

exec /vagrant/bin/proxy
