# backend/root/etc/init/proxy.conf.erb
start on local-filesystems
respawn

expect daemon

<% vars.each do |name, value| %>
env <%= name %>='<%= value %>'
<% end %>

exec <%= home_dir %>/proxy/backend/bin/proxy-backend
