# backend/root/etc/systemd/system/proxy.service.erb
[Unit]
Description=Proxy
After=network.target
StartLimitIntervalSec=0

[Service]
Type=forking
Restart=always
ExecStart=<%= home_dir %>/proxy/backend/bin/proxy-backend

<% vars.each do |name, value| %>
Environment="<%= name %>=<%= value %>"
<% end %>

[Install]
WantedBy=multi-user.target
